# workerbee/systemd/sync-datasets-to-mcs.service
[Unit]
Description=Sync datasets to MCS.

[Service]
ExecStart=/home/sean/sync-datasets.sh
