# workerbee/systemd/update-digests.timer
[Unit]
Description=Update digests nightly.

[Timer]
OnCalendar=*-*-* 19:30:00

[Install]
WantedBy=timers.target
