# workerbee/systemd/update-recent-digests.timer
[Unit]
Description=Update datasets hourly.

[Timer]
OnCalendar=*:02/5

[Install]
WantedBy=timers.target
