# workerbee/systemd/update-plenario-json.timer
[Unit]
Description=Update datasets hourly.

[Timer]
OnCalendar=*:04/5

[Install]
WantedBy=timers.target
