# workerbee/systemd/update-digests.service
[Unit]
Description=Update project digests.

[Service]
ExecStart=/home/sean/update-digests.sh
