# data-exporter/etc/systemd/system/refresh-recent-datasets.service
[Unit]
Description=Refresh datasets

[Service]
WorkingDirectory=/usr/lib/waggle/beehive-server/data-exporter
ExecStart=/usr/lib/waggle/beehive-server/data-exporter/export-recent-datasets-csv
