# data-exporter/etc/systemd/system/refresh-datasets.timer
[Unit]
Description=Refresh datasets hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
