# data-exporter/etc/systemd/system/refresh-recent-datasets.timer
[Unit]
Description=Refresh recent datasets

[Timer]
OnBootSec=15min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
