# data-exporter/etc/systemd/system/refresh-live-nodes.service
[Unit]
Description=Refresh live nodes

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/lib/waggle/beehive-server/bin/refresh-live-nodes'
