# data-exporter/etc/systemd/system/refresh-datasets.service
[Unit]
Description=Refresh datasets

[Service]
Type=oneshot
WorkingDirectory=/usr/lib/waggle/beehive-server/data-exporter
ExecStart=/usr/lib/waggle/beehive-server/data-exporter/refresh-datasets.sh
