# data-exporter/etc/systemd/system/sync-datasets-to-mcs.service
[Unit]
Description=Sync datasets to MCS

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'rsync -v -a /mnt/datasets/mcs/datasets/ mcswww:/mcs/www.mcs.anl.gov/research/projects/waggle/downloads/datasets'
