# workerbee-old/systemd/update-docs.service
[Unit]
Description=Update waggle docs

[Service]
ExecStart=/home/sean/crondoc/update
