# workerbee-old/systemd/update-recent-datasets.service
[Unit]
Description=Update decoded datasets from last 24h

[Service]
ExecStart=/home/sean/update-recent-datasets.sh
