# workerbee-old/systemd/update-docs.timer
[Unit]
Description=Update docs.

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
