# workerbee-old/systemd/update-datasets.timer
[Unit]
Description=Update datasets hourly.

[Timer]
OnCalendar=*:00:00

[Install]
WantedBy=timers.target
